# autostart/magic.desktop
[Desktop Entry]
Type=Application
Name=MagicMirror
Exec=/home/razvanh/MagicMirro/startw.sh