# autostart/mirror.desktop
[Desktop Entry]

Name=MagicMirror

Type=Application

Exec=/home/razvanh/Desktop/startw.sh