# autostart/mirror.desktop
[Desktop Entry]

Name=MagicMirror

Type=Application

Exec=/home/razvanh/MagicMirror/startw.sh
